# n.conf
[program:n]
command=python3 /home/curator2/n/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/n/test.err.log
stdout_logfile=/home/curator2/n/test.out.log